# low.slice
# ATL
[Slice]
#CPUAccounting=yes
CPUWeight=2
#IOAccounting=yes
IOWeight=2
#MemoryAccounting=yes
MemoryHigh=35%
MemoryMax=50%

[Unit]
Description=Low priority level
