# out/high.slice
[Unit]
Description=High priority level
Documentation=man:systemd.slice
Before=simple-slices.target

[Slice]
#CPUAccounting=yes
CPUWeight=200
#IOAccounting=yes
IOWeight=200
#MemoryAccounting=yes
MemoryLow=40%
MemoryHigh=90%
MemoryMax=99%

#[Install]
#WantedBy=simple-slices.target
