# out/low.slice




[Unit]
Description=Low priority level
Documentation=man:systemd.slice
Before=slices.target

[Slice]
#CPUAccounting=yes
CPUWeight=2
#IOAccounting=yes
IOWeight=2
#MemoryAccounting=yes

MemoryHigh=35%
MemoryMax=50%

#[Install]
#WantedBy=slices.target


