# out/mh1.slice
[Unit]
Description=Lowest mhp level
Documentation=man:systemd.slice
Before=simple-slices.target

[Slice]
#CPUAccounting=yes
CPUWeight=75
#IOAccounting=yes
IOWeight=75
#MemoryAccounting=yes
MemoryLow=20%
MemoryHigh=85%
MemoryMax=99%

#[Install]
#WantedBy=simple-slices.target
