# out/mh2.slice






[Unit]
Description=Highest mhp level
Documentation=man:systemd.slice
Before=slices.target

[Slice]
#CPUAccounting=yes
CPUWeight=125
#IOAccounting=yes
IOWeight=125
#MemoryAccounting=yes
MemoryLow=20%
MemoryHigh=85%
MemoryMax=99%

#[Install]
#WantedBy=slices.target


