# out/ml1.slice





[Unit]
Description=Lowest mlp level
Documentation=man:systemd.slice
Before=slices.target

[Slice]
#CPUAccounting=yes
CPUWeight=10
#IOAccounting=yes
IOWeight=10
#MemoryAccounting=yes

MemoryHigh=80%
MemoryMax=95%

#[Install]
#WantedBy=slices.target


