# out/simple-slices.target
[Unit]
Description=A small set of systemd slices with various priorities
Documentation=man:systemd.resource-control
Wants=low.slice ml1.slice ml2.slice mh1.slice mh2.slice high.slice

[Install]
WantedBy=slices.target
